# build/supervisord.conf
[supervisord]
nodaemon=true

[program:nscd]
command=/usr/sbin/nscd -F
autostart=true
autorestart=true

[program:threadfin]
command=sh -c "${THREADFIN_BIN}/threadfin -port=${THREADFIN_PORT} -bind=${THREADFIN_BIND_IP_ADDRESS} -config=${THREADFIN_CONF} -debug=${THREADFIN_DEBUG}"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
