# build/supervisord.conf
[supervisord]
nodaemon=true

[program:nscd]
command=/usr/sbin/nscd -F
autostart=true
autorestart=true

[program:threadfin]
command=%(ENV_THREADFIN_BIN)s/threadfin -port=%(ENV_THREADFIN_PORT)s -bind=%(ENV_THREADFIN_BIND_IP_ADDRESS)s -config=%(ENV_THREADFIN_CONF)s -debug=%(ENV_THREADFIN_DEBUG)s
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
